# backend/ignition_toolkit/clouddesigner/docker_files/designer-desktop/config/xfce4/autostart/designer.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Ignition Designer Autostart
Comment=Auto-launch Ignition Designer on login with credentials
Exec=/usr/local/bin/launch-designer.sh
Terminal=false
Hidden=false
X-GNOME-Autostart-enabled=true
X-XFCE-Autostart-Override=true
StartupNotify=false
